# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
autostart=true
command=/usr/sbin/sshd -D

[program:liked]
autostart=true
command=/usr/bin/liked start
